# unit/dbus.socket
[Unit]
Description=dbus session bus

[Socket]
ListenStream=%t/dbus/user_bus_socket
