# unit/gnome.target
[Unit]
Description=gnome desktop
Wants=gnome-shell.service pulseaudio.service

[Install]
WantedBy=default.target