# unit/lightSM@.service.in
[Unit]
Description=Graphical session for %I

[Service]
Type=notify
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/lightSM -vu %I

[Install]
WantedBy=graphical.target